# bge_m3.service
[Unit]
Description=BGE-M3 Embedding Deployment Service (Docker + sglang)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash /home/murphy/local-deployment/bge_m3.sh
WorkingDirectory=/home/murphy/local-deployment
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
