# glm_4.7_flash.service
[Unit]
Description=GLM-4.7-Flash Deployment Service (Docker + vLLM)
After=network-online.target docker.service
Wants=network-online.target
Wants=glm_4.7_flash_tunnel.service

[Service]
Type=simple
ExecStart=/bin/bash /home/murphy/local-deployment/glm_4.7_flash.sh
WorkingDirectory=/home/murphy/local-deployment
Restart=always
RestartSec=10
Environment=ENABLE_SSH_TUNNEL=0

[Install]
WantedBy=default.target
