# qwen3_coder_30b.service
[Unit]
Description=Qwen3-Coder-30B Deployment Service
After=network-online.target
Wants=network-online.target
Wants=qwen3_coder_30b_tunnel.service

[Service]
Type=simple
ExecStart=/bin/bash /home/murphy/Local-Deployment-Scripts/qwen3_coder_30b.sh
WorkingDirectory=/home/murphy/Local-Deployment-Scripts
Restart=always
RestartSec=10
Environment=CUDA_VISIBLE_DEVICES=0
Environment=ENABLE_SSH_TUNNEL=0

[Install]
WantedBy=default.target
